# rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv sh

exec /usr/sbin/crond -f
# exec /usr/sbin/crond -c /etc/crontabs -L /var/log/crond.log -l 6 -f;